# deploy/systemd/aidaytrader-bot.service
[Unit]
Description=Aidaytrader Algorithmic Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=trader
WorkingDirectory=/opt/aidaytrading
ExecStart=/opt/aidaytrading/run_bot.sh
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
